# dot_config/systemd/user/mangowc.target
[Unit]
Description=MangoWC Session

BindsTo=graphical-session.target
StopWhenUnneeded=yes
