# dot_config/systemd/user/waybar.service
[Unit]
Description=Waybar Service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/waybar

[Install]
WantedBy=graphical-session.target
